# scheduler/systemd/buildathon-tracker.service
[Unit]
Description=Buildathon Radar tracker service (Track/Applied click endpoints)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/jcube/projects/buildathon-radar
ExecStart=/home/jcube/projects/buildathon-radar/venv/bin/uvicorn buildathon_radar.tracker_service:app --host 127.0.0.1 --port 8015
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
