# scheduler/systemd/buildathon-radar.service
[Unit]
Description=Buildathon Radar weekly AI-event digest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/jcube/projects/buildathon-radar
ExecStart=/home/jcube/projects/buildathon-radar/venv/bin/python main.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
